# packaging/systemd/sigdep-sync.service
[Unit]
Description=SIGDEP-3 Edge Sync Agent
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=sigdep
Group=sigdep
EnvironmentFile=/etc/sigdep-sync/sigdep-sync.env
ExecStart=/usr/bin/java -jar /opt/sigdep-sync/sigdep-sync.jar
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sigdep-sync

[Install]
WantedBy=multi-user.target
